# supervisord.conf
[supervisord]
nodaemon = true
# Logs es pujen a stackdriver, aixi que no es necessari
logfile = /dev/null
logfile_maxbytes = 0
pidfile = /var/run/supervisord.pid

[unix_http_server]
file = /var/run/supervisor.sock
chmod = 0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock

[program:worker]
command = php /app/worker.php
user = root
process_name = %(program_name)s_%(process_num)02d
numprocs = 2
autostart = true
autorestart = true
priority = 5
stopsignal = KILL
stdout_logfile = /app/worker.out.log
stdout_logfile_maxbytes = 0
stdout_logfile_backups = 0
stderr_logfile = /app/worker.err.log
stderr_logfile_maxbytes = 0
stderr_logfile_backups = 0
